# systemd/mdbot@.service
[Unit]
Description=MDBOT: %i

[Service]
ExecStart=/usr/bin/mdbot %i
