# systemd/mdbot@NYSE.timer
[Unit]
Description=MDBOT timer: NYSE

[Timer]
OnCalendar=Mon..Fri *-*-* 16:30:00 America/New_York
