# systemd/mdbot@NYSE.timer
[Unit]
Description=MDBOT timer: NYSE

[Timer]
OnActiveSec=0s
OnUnitActiveSec=10s
AccuracySec=0s
