# systemd/mdbot-test@XETRA.timer
[Unit]
Description=MDBOT timer: NYSE

[Timer]
OnCalendar=Mon..Fri *-*-* 17:45:00 Europe/Berlin

[Install]
WantedBy=timers.target

