# systemd/mdbot-test@LSE.timer
[Unit]
Description=MDBOT timer: NYSE

[Timer]
OnCalendar=Mon..Fri *-*-* 17:15:00 Europe/London

[Install]
WantedBy=timers.target

