# systemd/mdbot@EURONEXT.timer
[Unit]
Description=MDBOT timer: NYSE

[Timer]
OnCalendar=Mon..Fri *-*-* 17:45:00 Europe/Amsterdam

[Install]
WantedBy=timers.target

