# systemd/mdbot-test@NYSE.timer
[Unit]
Description=MDBOT timer: NYSE

[Timer]
OnCalendar=Mon..Fri *-*-* 16:15:00 America/New_York

[Install]
WantedBy=timers.target

